# src/autostart/auto.desktop
[Desktop Entry]
Type=Application
Exec=.
Terminal=true
Name=tohutoro
Comment=make an Autostart
X-GNOME-Autostart-enabled=true
